# deploy/assessment-agent-redeploy.service
[Unit]
Description=Assessment-Agent 야간 재배포 (main 기준)
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
# 배포 체크아웃 경로 — 설치 환경에 맞게 수정.
Environment=REPO_DIR=/opt/assessment-agent
WorkingDirectory=/opt/assessment-agent
ExecStart=/opt/assessment-agent/deploy/redeploy.sh
# docker 명령 권한이 있는 사용자로 실행 (docker 그룹 소속 사용자 또는 root).
User=root
