# deploy/assessment-agent-redeploy.timer
[Unit]
Description=Assessment-Agent 야간 재배포 타이머 (매일 03:00)

[Timer]
# 매일 03:00 (서버 로컬 타임존 기준). timedatectl 로 타임존 확인 권장.
OnCalendar=*-*-* 03:00:00
# 머신이 03:00 에 꺼져 있었으면 부팅 후 놓친 실행을 1 회 보충.
Persistent=true
# 동시 부하 분산용 지터(최대 5분).
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
